# bitcoind-insight/bitcoind/supervisord.conf
[supervisord]
nodaemon=true

[program:bitcoind]
command=/usr/bin/bitcoind
